# stacks/m700/conf.d/systemd/system/mnt-data.automount
[Unit]
Description=Automount /mnt/data on access

[Automount]
Where=/mnt/data
# unmount after 150s of inactivity
TimeoutIdleSec=150

[Install]
WantedBy=multi-user.target