# stacks/m700/conf.d/systemd/system/shutdown.service
[Unit]
Description=Daily shutdown at midnight

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl poweroff