# stacks/m700/conf.d/systemd/system/suspend.service
[Unit]
Description=Suspend the system

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl suspend
