# stacks/m700/conf.d/systemd/system/wake.service
[Unit]
Description=Wake the system

[Service]
Type=oneshot
ExecStart=/bin/true
