# stacks/m700/conf.d/systemd/system/wake.timer
[Unit]
Description=Daily wake at 08:00

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
WakeSystem=true

[Install]
WantedBy=timers.target
