# stacks/m700/conf.d/systemd/user/davis-backup-local.timer
[Unit]
Description=Daily Davis backup -> local

[Timer]
OnCalendar=*-*-* 12:22:00
Persistent=true
Unit=davis-backup@local.service

[Install]
WantedBy=timers.target