# stacks/m700/conf.d/systemd/user/davis-backup-onsite.timer
[Unit]
Description=Weekly Davis backup -> onsite (Fri 20:10)

[Timer]
OnCalendar=Fri 20:22:00
Persistent=true
Unit=davis-backup@onsite.service

[Install]
WantedBy=timers.target