# stacks/m700/conf.d/systemd/user/gitea-backup-offsite.timer
[Unit]
Description=Daily Gitea backup -> offsite

[Timer]
OnCalendar=Fri 22:10:00
Persistent=true
Unit=gitea-backup@offsite.service

[Install]
WantedBy=timers.target