# stacks/m700/conf.d/systemd/user/gitea-backup-onsite.timer
[Unit]
Description=Weekly Gitea backup -> onsite (Fri 20:10)

[Timer]
OnCalendar=Fri 20:10:00
Persistent=true
Unit=gitea-backup@onsite.service

[Install]
WantedBy=timers.target