# stacks/m700/conf.d/systemd/user/gitea-backup-nas.timer
[Unit]
Description=Weekly Gitea backup -> nas (Fri 20:10)

[Timer]
OnCalendar=Fri 20:10:00
Persistent=true
Unit=gitea-backup@nas.service

[Install]
WantedBy=timers.target