# stacks/m700/conf.d/systemd/user/lldap-certbot-renew.timer
[Unit]
Description=Weekly certbot certificate renewal check for lldap

[Timer]
OnCalendar=weekly
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target