# stacks/m700/conf.d/systemd/user/pocket-id-backup-offsite.timer
[Unit]
Description=Daily Pocket ID backup -> offsite

[Timer]
OnCalendar=Fri 22:14:00
Persistent=true
Unit=pocket-id-backup@offsite.service

[Install]
WantedBy=timers.target