# stacks/m700/conf.d/systemd/user/wallos-backup-offsite.timer
[Unit]
Description=Daily Wallos backup -> offsite

[Timer]
OnCalendar=Fri 22:18:00
Persistent=true
Unit=wallos-backup@offsite.service

[Install]
WantedBy=timers.target