# stacks/m700/conf.d/backup/systemd/backup-pocket-id.service
[Unit]
Description=Weekly backup for Pocket-ID

[Service]
Type=oneshot
WorkingDirectory=/home/max/home-ops/stacks/m700
ExecStart=/usr/local/bin/just backup-all pocket-id
