# stacks/m700/conf.d/systemd/backup-opengist.timer
[Unit]
Description=Run Opengist backup daily at 21:02

[Timer]
OnCalendar=Sun *-*-* 22:06:00
Persistent=true

[Install]
WantedBy=timers.target
